# deploy/systemd/rustmgrapi.service
[Unit]
Description=Rust Manager API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/rust-manager/api
ExecStart=/usr/bin/dotnet /opt/rust-manager/api/rustmgrapi.dll
Restart=always
RestartSec=5
User=rustmgr
Group=rustmgr
EnvironmentFile=-/opt/rust-manager/config.env
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=HOME=/opt/rust-manager
Environment=DOTNET_CLI_HOME=/opt/rust-manager
Environment=TMPDIR=/tmp

[Install]
WantedBy=multi-user.target
